# pihole-influx-exporter.cron
*/10 * * * * /usr/local/bin/pihole-influx-exporter >> /var/log/pihole-influx-exporter.log 2>&1
